# root/etc/s6-overlay/s6-rc.d/init-config-cuda/run
#!/usr/bin/with-contenv bash

if [ "${CUDA_ACCELERATION}" = "true" ]; then
    # install pip deps
    apt-get update
    apt-get install --no-install-recommends -y \
        g++ \
        python3-dev \
        python3-venv

    # create venv if it does not exist
    if [ ! -d /config/cuda-venv ]; then
        python3 -m venv /config/cuda-venv/
    fi

    # install/update pip packages to /config venv
    source /config/cuda-venv/bin/activate
    pip install -U --no-cache-dir --index-url https://download.pytorch.org/whl/cu117 \
        torch
    pip install -U --no-cache-dir \
        fastapi \
        insightface \
        onnxruntime \
        optimum \
        scikit-learn \
        scipy \
        sentence-transformers \
        transformers \
        uvicorn[standard]
    pip install --break-system-packages -U --no-cache-dir nvidia-pyindex && \
    pip install --break-system-packages -U --no-cache-dir onnx-graphsurgeon && \

    # cleanup
    apt-get remove -y --purge \
        g++ \
        python3-dev
    apt-get autoremove -y --purge
    apt-get clean
fi
