# root/etc/s6-overlay/s6-rc.d/svc-typesense/run
#!/usr/bin/with-contenv bash

# stop typesense from starting
if [ "${DISABLE_TYPESENSE}" == "true" ]; then
    tail -f /dev/null
fi

exec \
    cd /app/typesense s6-setuidgid abc \
        ./typesense-server \
            --api-key=${TYPESENSE_API_KEY} \
            --enable-cors &>/dev/null
