# root/etc/s6-overlay/s6-rc.d/svc-typesense/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/typesense



exec s6-setuidgid abc \
	./typesense-server --data-dir=$(pwd)/typesense-data --api-key=$TYPESENSE_API_KEY --enable-cors
