# root/etc/s6-overlay/s6-rc.d/svc-web/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/immich/web || exit 1

export CHOKIDAR_USEPOLLING=true

exec s6-setuidgid abc \
	node build/index.js
