# root/etc/s6-overlay/s6-rc.d/svc-web/run
#!/usr/bin/with-contenv bash

export HOME=/config

exec \
    cd /app/immich/web s6-setuidgid abc \
        node build/index.js
