# root/etc/s6-overlay/s6-rc.d/svc-microservices/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    cd /app/immich/server s6-setuidgid abc \
        node dist/main microservices || exit
