# root/etc/s6-overlay/s6-rc.d/init-test-run/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -euo pipefail

export DEBIAN_FRONTEND=noninteractive

curl -fsSL https://www.postgresql.org/media/keys/ACCC4CF8.asc | gpg --dearmor > /usr/share/keyrings/pgdg-archive-keyring.gpg && \
echo "deb [signed-by=/usr/share/keyrings/pgdg-archive-keyring.gpg] https://apt.postgresql.org/pub/repos/apt noble-pgdg main" > /etc/apt/sources.list.d/pgdg.list

wait_for_pg() {
    local retries=30
    while (( retries-- > 0 )); do
        if s6-setuidgid postgres pg_isready -q; then
            return 0
        fi
        sleep 1
    done
    echo "PostgreSQL did not become ready in time" >&2
    return 1
}

# install immich dependencies for test runs
if [[ -n "${TEST_RUN:-}" ]]; then
    echo "Configuring CI for test run"
    apt-get update &>/dev/null
    apt-get remove -y --purge \
        'postgresql*' &>/dev/null || true
    apt-get install -y --no-install-recommends \
        postgresql-16 \
        postgresql-16-pgvector \
        redis-server &>/dev/null

    if ! s6-setuidgid postgres pg_ctlcluster 16 main status &>/dev/null; then
        s6-setuidgid postgres pg_ctlcluster 16 main start
    fi
    wait_for_pg
    s6-setuidgid postgres psql -c "ALTER USER postgres WITH PASSWORD 'password';"
    s6-setuidgid postgres psql -c "CREATE EXTENSION vector;"

    if [[ "$(arch)" == "x86_64" ]]; then
        version="amd64"
    else
        version="arm64"
    fi

    # Install vchord
    curl -o \
        /tmp/vchord.deb -L \
        "https://github.com/tensorchord/VectorChord/releases/download/0.3.0/postgresql-16-vchord_0.3.0-1_$version.deb"
    dpkg -i /tmp/vchord.deb

    # Setup vchord
    s6-setuidgid postgres psql -c "ALTER EXTENSION vector UPDATE;"
    s6-setuidgid postgres psql -c 'ALTER SYSTEM SET shared_preload_libraries = "vchord.so"'
    s6-setuidgid postgres pg_ctlcluster 16 main restart
    wait_for_pg
    s6-setuidgid postgres psql -c "CREATE EXTENSION vchord CASCADE;"

    s6-setuidgid abc redis-server --dir /config/ &>/dev/null &

    echo "PostgreSQL/Redis started"
fi
