# root/etc/s6-overlay/s6-rc.d/init-test-run/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

curl -fsSL https://www.postgresql.org/media/keys/ACCC4CF8.asc | gpg --dearmor > /usr/share/keyrings/pgdg-archive-keyring.gpg && \
echo "deb [signed-by=/usr/share/keyrings/pgdg-archive-keyring.gpg] https://apt.postgresql.org/pub/repos/apt noble-pgdg main" > /etc/apt/sources.list.d/pgdg.list

# install immich dependencies for test runs
if [[ -n "${TEST_RUN}" ]]; then
    echo "Configuring CI for test run"
    apt-get update &>/dev/null
    apt remove -y --purge \
        postgresql* &>/dev/null
    apt-get install -y --no-install-recommends \
        postgresql-16 \
        postgresql-16-pgvector \
        redis-server &>/dev/null

    s6-setuidgid postgres pg_ctlcluster 16 main start
    s6-setuidgid postgres psql -c "ALTER USER postgres WITH PASSWORD 'password';"
    s6-setuidgid postgres psql -c "CREATE EXTENSION vector;"

    if [[ "$(arch)" == "x86_64" ]]; then
        version="amd64"
    else
        version="arm64"
    fi

    # Install vchord
    curl -o \
        /tmp/vchord.deb -L \
        "https://github.com/tensorchord/VectorChord/releases/download/0.3.0/postgresql-16-vchord_0.3.0-1_$version.deb"
    dpkg -i /tmp/vchord.deb

    # Setup vchord
    s6-setuidgid postgres psql -c "ALTER EXTENSION vector UPDATE;"
    s6-setuidgid postgres psql -c 'ALTER SYSTEM SET shared_preload_libraries = "vchord.so"'
    sleep 1
    pkill postgres
    sleep 1
    s6-setuidgid postgres pg_ctlcluster 16 main start
    sleep 1
    s6-setuidgid postgres psql -c "CREATE EXTENSION vchord CASCADE;"

    s6-setuidgid abc redis-server --dir /config/ &>/dev/null &

    echo "PostgreSQL/Redis started"
fi
