# root/etc/s6-overlay/s6-rc.d/svc-server/run
#!/usr/bin/with-contenv bash
if [[ -z "${TYPESENSE_HOST}" ]]; then
    export TYPESENSE_ENABLED="false"
fi

exec \
    cd /app/immich/server s6-setuidgid abc \
        node dist/main immich
