# root/etc/s6-overlay/s6-rc.d/svc-server/run
#!/usr/bin/with-contenv bash


export LD_PRELOAD="/usr/lib/$(arch)-linux-gnu/libmimalloc.so.2"

exec \
    s6-notifyoncheck -d -n 300 -w 5000 -c "nc -z localhost 8080" \
        cd /app/immich/server s6-setuidgid abc \
            node dist/main immich
