# root/etc/s6-overlay/s6-rc.d/init-config-immich/run
#!/usr/bin/with-contenv bash

# make folders
mkdir -p \
    /config/{typesense,log/nginx} \
    /photos \
    /run/immich

if [[ -z "${TEST_RUN}" ]]; then
    if [[ "$DOCKER_MODS" != *"postgres"* ]]; then
        if [ ! -f /config/.migrated ]; then
            echo "Attempting to automatically migrate the database"
            if python3 /defaults/migrate.py; then
                echo "Database migrated successfully. It is recommended to run a storage migration in Immich."
            else
                echo "ERROR: Database migration failed. Please investigate or open an issue."
                sleep infinity
            fi
        fi
    else
        if [ ! -f /config/.migrated ]; then
            echo "WARN: It looks like you're using the PostgreSQL docker mod, the database cannot automatically migrate."
            echo "WARN: The container will now sleep until you manually restart it"
            echo "WARN: You MUST run 'python3 /defaults/migrate.py' within the container when you start it again, otherwise all your photos won't appear"
            touch /config/.migrated
        fi
    fi
fi

# permissions
find /app/immich -path "*/node_modules" -prune -o -exec chown abc:abc {} +
chown -R abc:abc \
    /config \
    /photos \
    /run/immich \
    /app/typesense \
    /cache
