# root/etc/s6-overlay/s6-rc.d/init-config-immich/run
#!/usr/bin/with-contenv bash

# make folders
mkdir -p \
    /photos

# check if /photos is empty
if [ -z "$(ls -A /photos)" ]; then
    lsiown -R abc:abc /photos
fi

# permissions
find /app/immich -path "*/node_modules" -prune -o -exec chown abc:abc {} +
lsiown -R abc:abc \
    /config \
    /usr/src/resources
