# root/etc/s6-overlay/s6-rc.d/init-config-immich/run
#!/usr/bin/with-contenv bash

# make folders
mkdir -p \
    /config/{typesense,log/nginx} \
    /photos \
    /run/immich

if [[ -z "${TEST_RUN}" ]]; then
    if [[ "$DOCKER_MODS" != *"postgres"* ]]; then
        if [ ! -f /config/.migrated ]; then
            echo "Attempting to automatically migrate the database..."
            if python3 /defaults/migrate.py; then
                echo "Database migration completed successfully."
                echo "Please run a storage migration in Immich to complete the setup."
            else
                echo "ERROR: Database migration failed. Please check the logs or open an issue."
                sleep infinity
            fi
        fi
    else
        if [ ! -f /config/.migrated ]; then
            echo "WARNING: The PostgreSQL docker mod is being used, and the database cannot be migrated automatically."
            echo "The container will now sleep until you manually restart it."
            echo "To complete the setup, please run 'python3 /defaults/migrate.py' within the container when you start it again."
            echo "Otherwise, your photos may not appear correctly."
            touch /config/.migrated
            sleep infinity
        fi
    fi
fi

# permissions
find /app/immich -path "*/node_modules" -prune -o -exec chown abc:abc {} +
chown -R abc:abc \
    /config \
    /photos \
    /run/immich \
    /app/typesense \
    /cache
