# root/etc/s6-overlay/s6-rc.d/init-config-immich/run
#!/usr/bin/with-contenv bash

# make folders
mkdir -p /import

# permissions
find /photos -maxdepth 1 -type d -exec lsiown abc:abc {} \;
find /app/immich -path "*/node_modules" -prune -o -exec chown abc:abc {} +
lsiown -R abc:abc \
    /config \
    /usr/src/resources
