# root/etc/s6-overlay/s6-rc.d/init-config-immich/run
#!/usr/bin/with-contenv bash

# make folders
mkdir -p \
    /photos

# permissions
find /photos \( ! -user abc -o ! -group abc \) -exec lsiown abc:abc {} \;
find /app/immich -path "*/node_modules" -prune -o -exec chown abc:abc {} +
lsiown -R abc:abc \
    /config \
    /usr/src/resources
