# root/etc/s6-overlay/s6-rc.d/init-config-immich/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# hmm... yes. dont say anything, it works
HOST="127.0.0.1       immich-server immich-microservices immich-machine-learning immich-web"
if ! grep -q "$HOST" /etc/hosts; then
	echo "$HOST" >>/etc/hosts
fi

# make folders
mkdir -p \
	/config/log/nginx \
	/run/immich

# permissions
find /app/immich -path "*/node_modules" -prune -o -exec chown abc:abc {} +
chown -R abc:abc \
	/config \
	/photos \
	/run/immich
