# root/etc/s6-overlay/s6-rc.d/init-config-immich/run
#!/usr/bin/with-contenv bash

# make folders
mkdir -p \
    /config/{typesense,log/nginx,machine-learning} \
    /photos \
    /run/immich

# permissions
find /app/immich -path "*/node_modules" -prune -o -exec chown abc:abc {} +
find /config -path "cuda-venv" -prune -o -exec chown abc:abc {} +
chown -R abc:abc \
    /photos \
    /run/immich \
    /app/typesense
