# root/etc/s6-overlay/s6-rc.d/init-gpu-acceleration/up
/etc/s6-overlay/s6-rc.d/init-gpu-acceleration/run
