# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash

if [[ "${IMMICH_MACHINE_LEARNING_URL}" != "http://127.0.0.1:3003" ]]; then
    echo "WARN: The Machine-Learning service is disabled as 'IMMICH_MACHINE_LEARNING_URL' is set to '${IMMICH_MACHINE_LEARNING_URL}'."
fi

# stop machine learning from starting
if [[ "${DISABLE_MACHINE_LEARNING}" == "true" ]] || [[ "${IMMICH_MACHINE_LEARNING_URL}" != "http://127.0.0.1:3003" ]]; then
    tail -f /dev/null
fi

export HOME=/tmp
exec \
    cd /app/immich/machine-learning s6-setuidgid abc \
       /lsiopy/bin/serve run app.main:ingress --host 0.0.0.0 --port 3003
