# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash

exec \
    cd /app/immich/machine-learning s6-setuidgid abc \
        python3 src/main.py
