# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash

# stop machine learning from starting
if [ "${DISABLE_MACHINE_LEANRNING}" == "true" ]; then
    tail -f /dev/null
fi

export HOME=/tmp
source /lsiopy/.venv/bin/activate
exec \
    cd /app/immich/machine-learning s6-setuidgid abc \
        python3 src/main.py
