# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash


# stop machine learning from starting
if [ "${DISABLE_MACHINE_LEANRNING}" == "true" ]; then
    tail -f /dev/null
fi

export HOME=/config
exec \
    cd /app/immich/machine-learning s6-setuidgid abc \
        gunicorn src.main:server
