# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/immich/machine-learning || exit 1

exec s6-setuidgid abc \
	node dist/main.js
