# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash

# stop machine learning from starting
if [ "${DISABLE_MACHINE_LEANRNING}" == "true" ]; then
    tail -f /dev/null
fi

if [ "${CUDA_ACCELERATION}" = "true" ]; then
    source /config/cuda-venv/bin/activate
fi

export HOME=/tmp
exec \
    cd /app/immich/machine-learning s6-setuidgid abc \
       /lsiopy/bin/serve run app.main:ingress --host 0.0.0.0 --port 3003
