# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/immich/machine-learning || exit 1

export NODE_ENV=${NODE_ENV:-production}

exec s6-setuidgid abc \
	node dist/main.js >/config/log/immich/machine-learning.log
