# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash

if [ "${IMMICH_MACHINE_LEARNING_URL}" != "http://127.0.0.1:3003" ]; then
    echo "WARN: The Machine-Learning service is disabled as 'IMMICH_MACHINE_LEARNING_URL' is set to '${IMMICH_MACHINE_LEARNING_URL}'."
fi

# stop machine learning from starting
if [ "${DISABLE_MACHINE_LEARNING}" == "true" ] || [ "${IMMICH_MACHINE_LEARNING_URL}" != "http://127.0.0.1:3003" ]; then
    tail -f /dev/null
fi

export LD_PRELOAD="/usr/lib/$(arch)-linux-gnu/libmimalloc.so.2"
export HOME=/tmp

exec \
    cd /app/immich/machine-learning s6-setuidgid abc \
       gunicorn app.main:app \
        -k uvicorn.workers.UvicornWorker \
        -w ${MACHINE_LEARNING_WORKERS:=1} \
        -b 0.0.0.0:3003 \
        -t ${MACHINE_LEARNING_WORKER_TIMEOUT:=120} \
        --log-config-json log_conf.json
