# root/etc/s6-overlay/s6-rc.d/svc-machine-learning/run
#!/usr/bin/with-contenv bash

export HOME=/config

exec \
    cd /app/immich/machine-learning s6-setuidgid abc \
        gunicorn src.main:server
